# config/example.crontab
01 00 * * * cd /<abs-path-to-project-folder>/localmusic;./develop.sh nightly-tasks
00 08 * * 0 cd /<abs-path-to-project-folder>/localmusic;./develop.sh weekly-tasks
01 00 1,15 * * cd /<abs-path-to-project-folder>/localmusic;./develop.sh biweekly-tasks
